# deploy/webhook.service
[Unit]
Description=d4d-front deploy webhook server
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
# ↓↓ 실제 서버 계정/경로로 수정하세요 ↓↓
User=ubuntu
WorkingDirectory=/home/ubuntu/apps/d4d-front
ExecStart=/usr/bin/node /home/ubuntu/apps/d4d-front/deploy/webhook-server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
